# obsidian-scrap/com.user.obsidian-auto-claude.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.obsidian-auto-claude</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/saisiot/bin/run_claude_latest.sh</string>
    </array>

    <key>WatchPaths</key>
    <array>
        <string>/Users/saisiot/Desktop/SecondBrain/Clippings</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/saisiot/auto_claude.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/saisiot/auto_claude.log</string>

    <key>ThrottleInterval</key>
    <integer>1</integer>
</dict>
</plist>
