# fleet-note-taker/com.fleetnotetaker.memo2fleet.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.fleetnotetaker.memo2fleet</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/saisiot/code_workshop/fleet_note_taker_claude/watch_memo2fleet.sh</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/memo2fleet_watcher.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/memo2fleet_watcher_error.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/saisiot/code_workshop/fleet_note_taker_claude</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin:/opt/homebrew/bin</string>
    </dict>
</dict>
</plist>
